# updater/crontab
0 8 * * * python /usr/src/app/updater.py >> /var/log/cron.log 2>&1
